# .config/autostart/xrdb-read-Xresources.desktop
[Desktop Entry]
Name=xrdb-runner
Comment=Run initial xrdb to set Xresources
Type=Application
Exec=env bash -c 'xrdb -merge ~/.Xresources'
